# files/sysmon-mailing.timer
[Unit]
Description=Periodic SysMon user mailing
After=network.target mysql.service

[Timer]
OnCalendar=monthly

[Install]
WantedBy=timers.target
