# files/sysmon.timer
[Unit]
Description=Check all systems daily
After=network.target mysql.service

[Timer]
OnCalendar=*-*-* 4:00:00

[Install]
WantedBy=timers.target
